# supervisor_docker.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too


[program:cloud-app-iot]
directory=/usr/iot_cloud_apps/iot
command=python3 app.py
priority=1
autostart=true
autorestart=true
user=frappe


[program:cloud-app-iot-ws-api]
directory=/usr/iot_cloud_apps/iot_ws_api
command=python3 app.py
priority=1
autostart=true
autorestart=true
user=frappe


[program:cloud-app-iot-bridge]
directory=/usr/iot_cloud_apps/iot_bridge
command=python3 app.py
priority=1
autostart=true
autorestart=true
user=frappe


[group:cloud-apps]
programs=cloud-app-iot,cloud-app-iot-ws-api,cloud-app-iot-bridge
